# Documents/MyMacSpace123/logging_project/github_code/loges-services-deploy-master/deploy/roles/logstash_objstore/files/monitorCronTab.cron
*/10 * * * * /opt/monitorKafka.sh >/dev/null 2>&1
