# Documents/MyMacSpace123/logging_project/github_code/loges-services-deploy-master/deploy/roles/monitor/templates/etc/supervisor/conf.d/monitor-backend.conf
[program:monitor-backend]
command=/usr/local/bin/runMonitorBackend.sh
stdout_logfile=/var/log/monitor-backend.log
redirect_stderr=true
priority=3
stopasgroup=true
