# Documents/MyMacSpace123/logging_project/github_code/loges-services-deploy-master/deploy/roles/monitor/templates/etc/supervisor/conf.d/monitor-collectd.conf
[program:monitor-collectd]
command=/usr/local/bin/runMonitorCollectd.sh
stdout_logfile=/var/log/monitorCollectd.log
redirect_stderr=true
priority=2
stopasgroup=true
