# Documents/MyMacSpace123/logging_project/github_code/loges-services-deploy-master/deploy/roles/monitor/templates/etc/supervisor/conf.d/monitor-frontend.conf
[program:monitor]
command=/usr/local/bin/runMonitorFrontend.sh
stdout_logfile=/var/log/monitor-frontend.log
redirect_stderr=true
priority=4
stopasgroup=true
