# Documents/MyMacSpace123/logging_project/github_code/loges-services-deploy-master/deploy/roles/monitor/templates/etc/supervisor/conf.d/provisioning-tests.conf
[program:provision-logmet]
command=/usr/local/bin/runLogmetProvision.sh